# supervisord.conf
[supervisord]
nodaemon=true

[program:start-script]
command=/usr/local/bin/start.sh
autostart=true
autorestart=false
startsecs=0
priority=1

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm/php-fpm8.1.err.log
stdout_logfile=/var/log/php-fpm/php-fpm8.1.out.log
startretries=3
priority=10

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/nginx.err.log
stdout_logfile=/var/log/nginx/nginx.out.log
startretries=3
priority=20
